# src/etc/systemd/system/auto-pa.service
[Unit]
Description=Node.js daemon for auto PA

[Service]
# TODO:
ExecStart=